# vscode/rootfs/etc/services.d/code/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Visual Studio Code
# Runs the code server
# ==============================================================================
declare -a options

bashio::log.info 'Starting the code server...'

# Non-interactive
options+=(--port 1337)
options+=(--data-dir "/data/vscode")

if bashio::config.false 'ssl'; then
    options+=(--allow-http)
else
    options+=(--cert "/ssl/$(bashio::config 'certfile')")
    options+=(--cert-key "/ssl/$(bashio::config 'keyfile')")
fi

if ! bashio::config.has_value 'password'; then
    options+=(--no-auth)
else
    options+=(--password "$(bashio::config 'password')")
fi

# Export env variables for the Home Assistant extension
export HASS_SERVER="http://hassio/homeassistant"
export HASS_TOKEN="${HASSIO_TOKEN:-}"

# Run the code server
exec code-server "${options[@]}" /config /share /adons /backups
